# deploy/airwave-update.path
[Unit]
Description=Watch for Airwave update requests

[Path]
PathExists=/var/lib/airwave/update.request
Unit=airwave-update.service

[Install]
WantedBy=multi-user.target
