# deploy/airwave-update.service
[Unit]
Description=Install an Airwave release update
After=network-online.target
Wants=network-online.target

[Service]
Type=oneshot
EnvironmentFile=-/etc/airwave.env
ExecStart=/usr/local/libexec/airwave-update
NoNewPrivileges=true
PrivateTmp=true
ProtectHome=true
ProtectSystem=strict
ReadWritePaths=/usr/local/bin /var/lib/airwave
